# deploy/systemd/nifty-auto-paper.service
[Unit]
Description=NIFTY Auto Paper Runtime Supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/nifty
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ubuntu/nifty/scripts/auto_paper_runtime.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
